# etc/zfs-auto-snapshot.cron.monthly
#!/bin/sh
exec flock -n /run/lock/zfs-auto-snapshot.lck zfs-auto-snapshot --quiet --syslog --label=monthly --keep=6 --pre-snapshot=/etc/zfs-auto-snapshot/pre //
